# vendor/home/python/home-daemon_old_20180203/home.service
[Unit]
Description=Home service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /var/www/html/home/vendor/LowPowerLab/home-daemon/hello_world.py
Restart=on-abort

[Install]
WantedBy=multi-user.target